# backend/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
loglevel=info
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/php-fpm-error.log
stderr_logfile_maxbytes=10MB
priority=1

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx/access.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/nginx/error.log
stderr_logfile_maxbytes=10MB
priority=2
depends_on=php-fpm

[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stdout_logfile=/var/log/queue-worker.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/queue-worker-error.log
stderr_logfile_maxbytes=10MB
priority=3
stopwaitsecs=3600
